# roles/networkmanager/templates/dispatcher.d/50-vpn-mount-smb.sh.j2
#jinja2:lstrip_blocks: True
#!/bin/sh
{{ ansible_managed | comment }}

# This script mounts SMB shares associated with a specific VPN connection
# when the VPN connection comes up.
# Each share is mounted per user into their home directory.

# NM_VPN_NAME is provided by NetworkManager dispatcher for VPN events.
LOG_PREFIX="NM-VPN-SMB-Mount"

{% set id_map = {} -%}
{% if network_connections_data is defined -%}
  {% for conn_type in ['ethernet', 'wifi'] -%}
    {% if network_connections_data[conn_type] is defined -%}
      {% for conn_name, conn_details in network_connections_data[conn_type].items() -%}
        {% if conn_details.uuid is defined -%}
          {% set _ = id_map.update({conn_name: conn_details.uuid}) -%}
        {% endif -%}
      {% endfor -%}
    {% endif -%}
  {% endfor -%}
  {% if network_connections_data.vpn is defined -%}
    {% for vpn_name, vpn_details in network_connections_data.vpn.items() -%}
      {% set _ = id_map.update({vpn_name: vpn_name}) -%}
    {% endfor -%}
  {% endif -%}
{% endif -%}
{% set user_homes = {} -%}
{% for user in users_list | default([]) -%}
  {% set _ = user_homes.update({user.name: user.home | default('/home/' + user.name)}) -%}
{% endfor -%}
if [ "$2" = "vpn-up" ]; then
{% if networkmanager_smb_shares is defined and networkmanager_smb_shares is not none and networkmanager_smb_shares | length > 0 %}
  {% for share in networkmanager_smb_shares %}
    {% set dependency_vpn_names_for_this_share = [] %}

    {% for dependency_name in share.depends_on %}
      {% if id_map[dependency_name] is defined and id_map[dependency_name] == dependency_name %}
        {% set _ = dependency_vpn_names_for_this_share.append(id_map[dependency_name]) %}
      {% endif %}
    {% endfor %}

    {% if dependency_vpn_names_for_this_share | length > 0 %}
  # Check if the current NM_VPN_NAME matches any of the share's VPN connection dependencies
  if {% for vpn_name in dependency_vpn_names_for_this_share %}[ "$NM_VPN_NAME" = "{{ vpn_name }}" ]{% if not loop.last %} || {% endif %}{% endfor %}; then
      {% for smb_user in share.users %}
        {% set user_home = user_homes[smb_user.name] | default('/home/' + smb_user.name) %}
        {% set mount_point = user_home + '/' + share.mount_subpath %}
        {% set creds_path = user_home + '/' + networkmanager_smb_credentials_dir + '/' + share.credentials_name + '-credentials' %}
        {% set opts = ['uid=$(id -u ' + smb_user.name + ')', 'gid=$(id -g ' + smb_user.name + ')', 'credentials=' + creds_path] %}
        {% for key, value in (share.mount_options | default({})).items() %}
          {% set _ = opts.append(key + '=' + value | string) %}
        {% endfor %}
    # Mount for user {{ smb_user.name }}
    if ! mountpoint -q "{{ mount_point }}"; then
      echo "$LOG_PREFIX: Current VPN connection ($NM_VPN_NAME) is a dependency for share '{{ share.name }}'. Mounting {{ mount_point }} for {{ smb_user.name }}"
      mkdir -p "{{ mount_point }}"
      mount -t cifs -o {{ opts | join(',') }} {{ share.src }} "{{ mount_point }}"
    else
      echo "$LOG_PREFIX: Share '{{ share.name }}' ({{ mount_point }}) already mounted for {{ smb_user.name }}."
    fi
      {% endfor %}
  fi
    {% endif %}
  {% endfor %}
{% endif %}
fi
