# roles/pki/templates/pki-cert-check.service.j2
{{ ansible_managed | comment }}
#
# Certificate expiry check service

[Unit]
Description=PKI Certificate Expiry Check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pki-cert-check
